# supervisord.conf
[supervisord]
nodaemon=true

[program:waarp-gwftp]
command=/usr/bin/waarp-gwftp start