# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:waarp-r66]
command=/bin/bash -c "/docker-entrypoint.sh"